# deploy/systemd/giga_backend.service
[Unit]
Description=Gigaotvet backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/gigaotvet/backend
EnvironmentFile=/opt/gigaotvet/backend/.env
ExecStart=/usr/bin/env uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
