# deploy/systemd/giga_frontend.service
[Unit]
Description=Gigaotvet frontend (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/gigaotvet/frontend
EnvironmentFile=/opt/gigaotvet/frontend/.env.production
ExecStart=/usr/bin/env next start -p 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
